# build/sha3sum/etc/systemd/system/sha3sum.service
[Unit]
Description=Tool sha3sum

[Service]
Type=oneshot
WorkingDirectory=/
ExecStart=/usr/bin/sha3sum

[Install]
WantedBy=multi-user.target

